# debian/motion.upstart
# Motion Daemon 

start on net-device-up
stop on runlevel [06]

description "Motion Detection"

respawn

exec /usr/bin/motion

